# files/seedlink.service
[Unit]
Description=Seedlink server

[Service]
User=sysop

ExecStart=/opt/seiscomp/bin/seiscomp start seedlink
ExecStop=/opt/seiscomp/bin/seiscomp stop seedlink

[Install]
WantedBy=multi-user.target
